# sandbox-runner/deploy/sandbox-server.service
# M8: production judge sandbox, wired into judge.ts/route.ts via
# src/lib/execute.ts (see itouOJ's own repo) for c/c++/python/javascript.
# Java still goes through Piston until sandbox-runner grows a Java profile.
# Runs as root: jail itself needs CAP_SYS_ADMIN-class privilege to create
# namespaces/cgroups (see sandbox-runner/src/jail.c and the plan's
# privilege-model notes) -- sandbox-server just fork/execs jail as a
# subprocess per request, so it needs to already have that privilege to
# hand down.
[Unit]
Description=Sandbox judge HTTP service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sandbox-runner
ExecStart=/opt/sandbox-runner/sandbox-server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
