# deploy/online-judge.service
[Unit]
Description=Online Judge (Next.js)
After=network.target docker.service

[Service]
Type=simple
User=oj
Group=oj
WorkingDirectory=/opt/online-judge
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=3
Environment=NODE_ENV=production

# 基本強化
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
